# systemd/lobster-compress.service
[Unit]
Description=LobsterPress Compressor - 预测性压缩引擎
After=network.target

[Service]
Type=oneshot
ExecStart=%h/bin/predictive-compressor.sh
WorkingDirectory=%h
Environment=AUTO_APPLY=true
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=default.target
