# systemd/lobster-optimizer.service
[Unit]
Description=LobsterPress Optimizer - 成本优化器
After=network.target

[Service]
Type=oneshot
ExecStart=%h/bin/cost-optimizer.sh
WorkingDirectory=%h
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=default.target
